# root/etc/s6-overlay/s6-rc.d/init-altus/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
echo "**** Fixing app folder permissions ****"
lsiown -R abc:abc \
    /opt/altus
